# deploy/debian/webmanager-update.path
[Unit]
Description=Install a super-admin-approved WebManager update

[Path]
PathChanged=/var/lib/webmanager-updater/requests/install.commit
Unit=webmanager-update.service

[Install]
WantedBy=multi-user.target

