# deploy/debian/webmanager-update.path
[Unit]
Description=Run a requested WebManager update check or approved installation

[Path]
PathChanged=/var/lib/webmanager-updater/requests/install.commit
PathChanged=/var/lib/webmanager-updater/requests/check
Unit=webmanager-update.service

[Install]
WantedBy=multi-user.target
